# config/hourly.cron
5 * * * *       /srv/unsourced.org/unsourced/hourlyjobs >/dev/null
